# rootfs/etc/services.d/supervisor-gateway/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Supervisor API Gateway service
# ==============================================================================
bashio::log.info "Starting Supervisor API Gateway..."

# Print environment for debugging
bashio::log.info "Checking for SUPERVISOR_TOKEN..."
if [ -n "${SUPERVISOR_TOKEN}" ]; then
    bashio::log.info "SUPERVISOR_TOKEN found in environment"
else
    bashio::log.warning "SUPERVISOR_TOKEN not found in environment"
fi

# Start the Python application
cd /app || bashio::exit.nok "Cannot change to /app directory"
exec python3 /app/app.py
